# var/lib/dpkg/info/mariadb-server-10.0.list
/.
/usr
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_mariadb-10.0.py
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mariadb-server-10.0
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/myisamchk.1.gz
/usr/share/man/man1/aria_pack.1.gz
/usr/share/man/man1/mysql_convert_table_format.1.gz
/usr/share/man/man1/mysql_install_db.1.gz
/usr/share/man/man1/aria_ftdump.1.gz
/usr/share/man/man1/resolve_stack_dump.1.gz
/usr/share/man/man1/mysql_plugin.1.gz
/usr/share/man/man1/myisampack.1.gz
/usr/share/man/man1/myisamlog.1.gz
/usr/share/man/man1/replace.1.gz
/usr/share/man/man1/mysqld_multi.1.gz
/usr/share/man/man1/mysql_zap.1.gz
/usr/share/man/man1/mysql_secure_installation.1.gz
/usr/share/man/man1/mysql_setpermission.1.gz
/usr/share/man/man1/msql2mysql.1.gz
/usr/share/man/man1/aria_dump_log.1.gz
/usr/share/man/man1/mysqlhotcopy.1.gz
/usr/share/man/man1/my_print_defaults.1.gz
/usr/share/man/man1/resolveip.1.gz
/usr/share/man/man1/perror.1.gz
/usr/share/man/man1/myisam_ftdump.1.gz
/usr/share/man/man1/aria_read_log.1.gz
/usr/share/man/man1/aria_chk.1.gz
/usr/share/man/man1/mysqld_safe.1.gz
/usr/share/man/man1/mysqlbinlog.1.gz
/usr/share/man/man1/mysql_tzinfo_to_sql.1.gz
/usr/share/mysql
/usr/share/mysql/fill_help_tables.sql
/usr/share/mysql/mysql_performance_tables.sql
/usr/share/mysql/errmsg-utf8.txt
/usr/share/mysql/debian-start.inc.sh
/usr/share/mysql/mysql_system_tables_data.sql
/usr/share/mysql/echo_stderr
/usr/share/mysql/mysql_test_data_timezone.sql
/usr/share/mysql/mysql_system_tables.sql
/usr/share/doc
/usr/share/doc/mariadb-server-10.0
/usr/share/doc/mariadb-server-10.0/copyright
/usr/share/doc/mariadb-server-10.0/mysqld.sym.gz
/usr/share/doc/mariadb-server-10.0/README.Debian.gz
/usr/bin
/usr/bin/replace
/usr/bin/aria_dump_log
/usr/bin/aria_pack
/usr/bin/mysqld_safe
/usr/bin/myisam_ftdump
/usr/bin/aria_ftdump
/usr/bin/mysql_tzinfo_to_sql
/usr/bin/mysqld_safe_helper
/usr/bin/msql2mysql
/usr/bin/aria_chk
/usr/bin/mysql_convert_table_format
/usr/bin/myisamlog
/usr/bin/perror
/usr/bin/mysql_zap
/usr/bin/myisamchk
/usr/bin/mysqld_multi
/usr/bin/resolve_stack_dump
/usr/bin/mysql_plugin
/usr/bin/mysql_setpermission
/usr/bin/mysqlhotcopy
/usr/bin/resolveip
/usr/bin/mysql_install_db
/usr/bin/myisampack
/usr/bin/my_print_defaults
/usr/bin/aria_read_log
/usr/bin/mysql_secure_installation
/usr/bin/mysqlbinlog
/usr/lib
/usr/lib/mysql
/usr/lib/mysql/plugin
/usr/lib/mysql/plugin/metadata_lock_info.so
/usr/lib/mysql/plugin/auth_pam.so
/usr/lib/mysql/plugin/handlersocket.so
/usr/lib/mysql/plugin/sql_errlog.so
/usr/lib/mysql/plugin/ha_sequence.so
/usr/lib/mysql/plugin/ha_innodb.so
/usr/lib/mysql/plugin/ha_sphinx.so
/usr/lib/mysql/plugin/semisync_slave.so
/usr/lib/mysql/plugin/semisync_master.so
/usr/lib/mysql/plugin/server_audit.so
/usr/lib/mysql/plugin/query_cache_info.so
/usr/lib/mysql/plugin/query_response_time.so
/usr/lib/mysql/plugin/locales.so
/etc
/etc/logrotate.d
/etc/logrotate.d/mysql-server
/etc/mysql
/etc/mysql/debian-start
/etc/mysql/mariadb.conf.d
/etc/mysql/mariadb.conf.d/50-mysqld_safe.cnf
/etc/mysql/mariadb.conf.d/50-server.cnf
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.mysqld
/etc/default
/etc/default/mysql
/etc/logcheck
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/mariadb-server-10_0
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/mariadb-server-10_0
/etc/logcheck/ignore.d.workstation
/etc/logcheck/ignore.d.workstation/mariadb-server-10_0
/etc/init.d
/etc/init.d/mysql
/var
/var/lib
/var/lib/mysql-upgrade
/usr/share/doc/mariadb-server-10.0/changelog.Debian.gz
